# supervisord.conf
[program:wikinotify]
command = python /home/dotcloud/current/wikibot.py
